# inkstory.service
[Unit]
Description=InkStory Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Замени 'pi' на своего пользователя (команда: whoami)
User=pi
WorkingDirectory=/home/pi/inkstory
ExecStart=/home/pi/inkstory/venv/bin/python bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
